# user-statistics-service/supervisord.conf
[supervisord]
nodaemon=true

[program:user-statistics-service]
command=java -jar /data/user-statistics-service.jar server /data/config.yml
